# deployment/systemd/mission-control-notification.service
[Unit]
Description=OpenClaw Mission Control Notification Daemon
After=network.target

[Service]
Type=simple
User=nisheeth
WorkingDirectory=/home/nisheeth/.openclaw/workspace
ExecStart=/home/nisheeth/.bun/bin/bun /home/nisheeth/.openclaw/workspace/scripts/notification-daemon.ts
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
